# openfaito/files/etc/init.d/strongswan
#!/bin/sh /etc/rc.common
# FAITO,2014

START=72
BIN=ipsec
DEFAULT=/usr/sbin/$BIN

initconf() {
	echo "# /etc/ipsec.secrets - Auto generated by /etc/init.d/strongswan" > /etc/ipsec.secrets
	echo ""             >> /etc/ipsec.secrets
	echo " : PSK $(uci -q get strongswan.ipsec_secrets.psk)" >> /etc/ipsec.secrets
	
	echo "# /etc/ipsec.conf - Auto generated by /etc/init.d/strongswan" > /etc/ipsec.conf
	echo "config setup" >> /etc/ipsec.conf
	echo ""             >> /etc/ipsec.conf
	echo "conn L2TP-IPSEC" >> /etc/ipsec.conf
	echo "      ikelifetime=$(uci get -q strongswan.ipsec_config_L2TP_IPSEC.ikelifetime)" >> /etc/ipsec.conf
	echo "      keylife=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.keylife)" >> /etc/ipsec.conf
	echo "      rekeymargin=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.rekeymargin)" >> /etc/ipsec.conf
	echo "      keyingtries=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.keyingtries)" >> /etc/ipsec.conf
	echo "      keyexchange=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.keyexchange)" >> /etc/ipsec.conf
	echo "      authby=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.authby)" >> /etc/ipsec.conf
	echo "      left=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.left)" >> /etc/ipsec.conf
	echo "      leftfirewall=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.leftfirewall)" >> /etc/ipsec.conf
	echo "      auto=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.auto)" >> /etc/ipsec.conf	
	echo "      right=$(uci -q get strongswan.ipsec_config_L2TP_IPSEC.right)" >> /etc/ipsec.conf
}
	
start() {
	initconf
        [ "$(uci -q get strongswan.main.enable)" == "1" ]&& $DEFAULT start
}

stop() {
	$DEFAULT stop
}
restart() {
	$DEFAULT stop
	initconf
        [ "$(uci -q get strongswan.main.enable)" == "1"  ]&& $DEFAULT start 
}
                  
reload() {
	initconf
	$DEFAULT update
}
